# scripts/com.edgeboard.train.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>                <string>com.edgeboard.train</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>bash</string>
    <string>-lc</string>
    <string>export NVM_DIR="$HOME/.nvm" &amp;&amp; [ -s "$NVM_DIR/nvm.sh" ] &amp;&amp; . "$NVM_DIR/nvm.sh" &amp;&amp; cd PROJECT_DIR &amp;&amp; NODE_OPTIONS="--max-old-space-size=12288" npx tsx scripts/daily-cycle.ts &gt;&gt; data/training/meta/launchd.log 2&gt;&amp;1; npx tsx scripts/check-models.ts &gt;&gt; data/training/meta/launchd.log 2&gt;&amp;1</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>    <integer>3</integer>
    <key>Minute</key>  <integer>15</integer>
  </dict>
  <key>RunAtLoad</key>             <false/>
  <key>KeepAlive</key>             <false/>
  <key>StandardOutPath</key>       <string>PROJECT_DIR/data/training/meta/launchd-stdout.log</string>
  <key>StandardErrorPath</key>     <string>PROJECT_DIR/data/training/meta/launchd-stderr.log</string>
</dict>
</plist>
